# portainer/rootfs/etc/s6-overlay/s6-rc.d/portainer/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Portainer service
# ==============================================================================

bashio::log.info "Starting Portainer..."

cd /opt/portainer || bashio::exit.nok "Could not change to Portainer directory"

exec ./portainer \
    --data /data \
    --no-analytics
