# n8n/rootfs/etc/s6-overlay/s6-rc.d/n8n/run
#!/command/with-contenv bashio
# ==============================================================================
# Start n8n service
# ==============================================================================

bashio::log.info "Starting n8n..."

# Get configuration options
TIMEZONE=$(bashio::config 'timezone')
WEBHOOK_URL=$(bashio::config 'webhook_url')
ENCRYPTION_KEY=$(bashio::config 'encryption_key')

# Set timezone if provided
if [[ -n "${TIMEZONE}" ]]; then
    export GENERIC_TIMEZONE="${TIMEZONE}"
    export TZ="${TIMEZONE}"
    bashio::log.info "Setting timezone to ${TIMEZONE}"
fi

# Set webhook URL if provided
if [[ -n "${WEBHOOK_URL}" ]]; then
    export WEBHOOK_URL="${WEBHOOK_URL}"
    bashio::log.info "Setting webhook URL to ${WEBHOOK_URL}"
fi

# Set encryption key if provided
if [[ -n "${ENCRYPTION_KEY}" ]]; then
    export N8N_ENCRYPTION_KEY="${ENCRYPTION_KEY}"
    bashio::log.info "Encryption key configured"
fi

# Set n8n environment variables
export N8N_HOST="0.0.0.0"
export N8N_PORT="5678"
export N8N_PROTOCOL="http"
export N8N_USER_FOLDER="/data"
export N8N_CONFIG_FILES="/data"

# Create data directory structure
mkdir -p /data/.n8n
chown -R root:root /data

# Start n8n
bashio::log.info "Starting n8n workflow automation..."
exec n8n start
