# qbittorrent/rootfs/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bash
/init &
cd /opt/flood
node server/bin/www --port 3000 --rundir /config &
wait
