# portainer/rootfs/etc/services.d/portainer/run
#!/usr/bin/with-contenv bashio

# Start Portainer
exec /portainer
