# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB service
# ==============================================================================

# Read options from Home Assistant add-on config
ROOT_PWD=$(bashio::config 'root_password')
DB_NAME=$(bashio::config 'database')
DB_USER=$(bashio::config 'user')
DB_PASS=$(bashio::config 'user_password')
REPL_IP=$(bashio::config 'replicate_from_ip')
REPL_USER=$(bashio::config 'replicate_user')
REPL_PASS=$(bashio::config 'replicate_password')
NODE_NAME=$(bashio::config 'node_name')
NODE_IP=$(bashio::config 'node_address')
CLUSTER_NAME=$(bashio::config 'cluster_name')

bashio::log.info "Configuring MariaDB Galera cluster..."

# Create Galera configuration
mkdir -p /etc/mysql/mariadb.conf.d
cat > /etc/mysql/mariadb.conf.d/60-galera.cnf << EOF
[galera]
wsrep_on=ON
wsrep_provider=/usr/lib/galera/libgalera_smm.so
wsrep_cluster_address="gcomm://${REPL_IP}"
wsrep_cluster_name="${CLUSTER_NAME}"
wsrep_node_address="${NODE_IP}"
wsrep_node_name="${NODE_NAME}"
wsrep_sst_method=rsync
wsrep_sst_auth="${REPL_USER}:${REPL_PASS}"
binlog_format=row
default_storage_engine=InnoDB
innodb_autoinc_lock_mode=2
bind-address=0.0.0.0
EOF

# Initialize MariaDB if not already done
if [ ! -d "/var/lib/mysql/mysql" ]; then
    bashio::log.info "Initializing MariaDB database..."
    mysql_install_db --user=mysql --datadir=/var/lib/mysql --auth-root-authentication-method=normal
fi

bashio::log.info "Starting MariaDB with Galera cluster support..."

# Start MariaDB directly
exec mysqld_safe --wsrep-new-cluster --user=mysql --datadir=/var/lib/mysql
