# mariadb/rootfs/etc/s6-overlay/s6-rc.d/mariadb/run
#!/command/with-contenv bashio
# ==============================================================================
# Start MariaDB service
# ==============================================================================

bashio::log.info "Starting MariaDB with Galera cluster support..."

exec /run.sh
