# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/qbittorrent/run
#!/command/with-contenv bashio
# Home Assistant Community Add-on: qBittorrent
# Runs qBittorrent

bashio::log.info 'Starting qBittorrent...'

exec /init
