# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/qbittorrent/run
#!/command/with-contenv bashio
# ==============================================================================
# Start qBittorrent service
# ==============================================================================

bashio::log.info "Starting qBittorrent..."

# Create default configuration if it doesn't exist
if [[ ! -f /config/qBittorrent/qBittorrent.conf ]]; then
    bashio::log.info "Creating default qBittorrent configuration..."
    mkdir -p /config/qbittorrent
    cat > /config/qbittorrent/qbittorrent.conf << EOL
[AutoRun]
enabled=false

[BitTorrent]
Session\DefaultSavePath=/downloads
Session\Port=6881

[Preferences]
WebUI\Address=*
WebUI\Port=8080
WebUI\Username=admin
WebUI\Password_ha1=@ByteArray(f6fdffe48c908deb0f4c3bd36c032e72)
EOL
    bashio::log.info "Default configuration created. Default WebUI credentials are admin/adminadmin"
fi

# Set permissions
chown -R root:root /config/qbittorrent

# Start qBittorrent
exec qbittorrent-nox \
    --webui-port=8080 \
    --profile=/config
