# mariadb/rootfs/etc/s6-overlay/s6-rc.d/init-mariadb/up
#!/command/with-contenv bashio
# ==============================================================================
# Initialize MariaDB database and users
# ==============================================================================

# Read options from Home Assistant add-on config
ROOT_PWD=$(bashio::config 'root_password')
DB_NAME=$(bashio::config 'database')
DB_USER=$(bashio::config 'user')
DB_PASS=$(bashio::config 'user_password')
REPL_USER=$(bashio::config 'replicate_user')
REPL_PASS=$(bashio::config 'replicate_password')

bashio::log.info "Waiting for MariaDB to be ready..."

# Wait for MariaDB to be ready
until mysqladmin ping --silent; do
  sleep 1
done

bashio::log.info "Initializing MariaDB database and users..."

# Set root password if provided
if [ -n "$ROOT_PWD" ]; then
    mysql -uroot <<EOF
ALTER USER 'root'@'localhost' IDENTIFIED BY '$ROOT_PWD';
FLUSH PRIVILEGES;
EOF
fi

# Create database and user if not exist
mysql -uroot ${ROOT_PWD:+-p"$ROOT_PWD"} <<EOF
CREATE DATABASE IF NOT EXISTS \`$DB_NAME\`;
CREATE USER IF NOT EXISTS '$DB_USER'@'%' IDENTIFIED BY '$DB_PASS';
GRANT ALL PRIVILEGES ON \`$DB_NAME\`.* TO '$DB_USER'@'%';
CREATE USER IF NOT EXISTS '$REPL_USER'@'%' IDENTIFIED BY '$REPL_PASS';
GRANT RELOAD, LOCK TABLES, PROCESS, REPLICATION CLIENT ON *.* TO '$REPL_USER'@'%';
FLUSH PRIVILEGES;
EOF

bashio::log.info "MariaDB initialization complete!"
